# American Airlines Wifi Extender/LaunchDaemons/com.Finch13.WifiUtilityHelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.Finch13.WifiUtilityHelper</string>
    <key>MachServices</key>
    <dict>
        <key>com.Finch13.WifiUtilityHelper</key>
        <true/>
    </dict>
    <key>BundleProgram</key>
    <string>WifiUtilityHelper</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
